# esantri-api.container
[Unit]
Description=Edaimi Backend API (Rootless)
After=network-online.target

[Container]
# Menggunakan image Node.js
Image=docker.io/library/node:20-alpine
WorkingDir=/app

# Mengeset nama container agar bisa dipanggil oleh cloudflared tunnel
ContainerName=esantri-api

# Menghubungkan container langsung ke jaringan Podman 'global_net'
Network=global_net

# Mount direktori kode aplikasi (dengan label keamanan SELinux ':Z')
Volume=/webku/appsku/esantri/backend-api:/app:Z

# Menjalankan backend API
Exec=node dist/main.js

# Menggunakan berkas environment host RHEL
EnvironmentFile=/webku/appsku/esantri/backend-api/.env

# === SEKSI HEALTH CHECK ===
HealthCmd=node --input-type=module -e "import http from 'http'; http.get('http://localhost:8080/api/v1/health', (r) => { if (r.statusCode === 200) process.exit(0); else process.exit(1); }).on('error', () => process.exit(1));"
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=10s

[Install]
WantedBy=default.target
